# OpenLinkHub.service
[Unit]
Description=Open source interface for iCUE LINK System Hub

[Service]
User=nikola
WorkingDirectory=/opt/OpenLinkHub
ExecStart=/opt/OpenLinkHub/OpenLinkHub
ExecReload=/bin/kill -s HUP $MAINPID
RestartSec=5

[Install]
WantedBy=multi-user.target